# sib-docker-build/supervisord/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
childlogdir = /var/log/supervisor

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket
loglevel=debug

[include]
files=/etc/supervisord.d/*.conf

[program:crond]
command=/usr/sbin/crond -n
autorestart=true

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
killasgroup=true
stopasgroup=true
autostart=true
autorestart=true

[program:mysql]
command=/usr/bin/mysqld_safe --basedir=/usr
killasgroup=true
stopasgroup=true
autostart=true
autorestart=true

[program:swamp-condor]
command=/opt/swamp/htcondor/sbin/condor_master -f
environment=CONDOR_CONFIG=/opt/swamp/htcondor/etc/condor_config
killasgroup=true
stopasgroup=true
autostart=true
autorestart=true
stopsignal=SIGQUIT

[program:swamp]
command=/usr/local/sbin/supervisord_start_swamp.bash
user=swa-daemon
directory=/opt/swamp/run
killasgroup=true
stopasgroup=true
autostart=true
autorestart=true