# build_docker_viewer/services/supervisor/supervisor.d/supercronic.conf
[program:supercronic]
command=/usr/local/bin/supercronic /usr/local/etc/crontab
priority=100
autostart=false
