# build_docker_viewer/services/supervisor/supervisor.d/tomcat.conf
[program:tomcat]
command=/opt/tomcat/bin/catalina.sh run
startsecs=5
priority=700
stopsignal=INT
autostart=false
