# build_docker_viewer/services/supervisor/supervisor.d/viewer-management-service.conf
[program:viewer-management-service]
command=/usr/local/libexec/viewer-management-service
stopwaitsecs=60
priority=999
autostart=true
